# debian/pohoda-digest/etc/cron.d/pohoda-digest
# Pohoda-Digest daily report
0 6 * * * root /usr/bin/pohoda-digest --quiet
